# debian/haloradio-frontend.cron.d
#
# Regular cron jobs for the haloradio-frontend package
#
*/5 * * * * radio nice /usr/share/haloradio/halostatcron.sh

